# systemd/littercam-prune.timer
[Unit]
Description=LitterCam daily retention prune

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
