# systemd/littercam-capture.service
[Unit]
Description=LitterCam capture service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/git/littercam
ExecStart=%h/git/littercam/.venv/bin/littercam-capture
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
