# systemd/littercam-prune.service
[Unit]
Description=LitterCam retention prune

[Service]
Type=oneshot
WorkingDirectory=%h/git/littercam
ExecStart=%h/git/littercam/.venv/bin/littercam-prune
